# deployment/argos-final.service.d/30-v1-source.conf
[Service]
# Drop-in: repoint argos-final at the v1 worktree (last pre-Mk-II commit) so :5173
# serves the stable v1 UI while the dev branch's Mk II / Carbon migration runs
# on :5174 via argos-dev.service.
#
# Empty ExecStart= clears the inherited value before re-declaring — without this,
# systemd would APPEND a second ExecStart and refuse to start the unit
# (per systemd.service(5) ExecStart= semantics).
#
# Render via scripts/ops/install-services.sh — substitutes the V1_PROJECT_DIR
# placeholder (default: /home/<user>/code/Argos-v1; override via env) and the
# NODE_BIN placeholder (resolved nvm-aware via the install script).
#
# Removal: rm this file + sudo systemctl daemon-reload + restart argos-final.
WorkingDirectory=__V1_PROJECT_DIR__
ExecStart=
ExecStart=__NODE_BIN__ --import tsx __V1_PROJECT_DIR__/scripts/ops/prod-server.ts
Environment="PORT=5173"
