# deployment/argos-final.service
[Unit]
Description=Argos SDR & Network Analysis Console
After=network.target gpsd.service earlyoom.service argos-startup.service
Wants=gpsd.service earlyoom.service

[Service]
Type=simple
User=__SETUP_USER__
Group=__SETUP_USER__
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__PROJECT_DIR__/.env
Environment="NODE_ENV=production"
Environment="NODE_OPTIONS=--max-old-space-size=1024"
# Pre-spawn tmux-0 on startup so the dashboard Terminal panel opens instantly.
Environment="ARGOS_TERMINAL_PRESPAWN=1"
# Custom entrypoint: adapter-node's build/index.js has no 'upgrade' listener, so
# WebSocket routes (e.g. /terminal-ws) would 404. scripts/ops/prod-server.ts
# imports the polka handler from build/ and wires the WS upgrade events.
ExecStart=__NODE_BIN__ --import tsx __PROJECT_DIR__/scripts/ops/prod-server.ts
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos
# NoNewPrivileges is intentionally NOT set: Argos control endpoints
# (openwebrx, bluehood, sparrow, kismet, etc.) must invoke `sudo systemctl`
# to start/stop their native peer services. NoNewPrivileges=true blocks
# sudo's setuid and breaks all such lifecycle endpoints. Privilege is
# instead scoped via narrow NOPASSWD rules in /etc/sudoers.d/argos and
# /etc/sudoers.d/argos-openwebrx.
PrivateTmp=true

[Install]
WantedBy=multi-user.target
