# deployment/argos-dev-monitor.service
[Unit]
Description=Argos Development Monitor (Vite/Chromium Keepalive)
After=network.target

[Service]
Type=simple
User=jetson2
ExecStart=/bin/bash __PROJECT_DIR__/scripts/ops/keepalive-dev.sh
WorkingDirectory=__PROJECT_DIR__
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
